# docker/dev/supervisord/sshd.conf
[program:sshd]
directory=/home/developer
command=/usr/sbin/sshd -D